# docker/8.1/cron/crontab
* * * * * root /usr/local/bin/php /var/www/html/app/artisan test:command
